# Assignment04/99-usb-kbd.rules
SUBSYSTEM=="input", ACTION=="add", SUBSYSTEMS=="usb", RUN+="/sbin/modprobe main"
